# templates/kopia.timer.j2
[Unit]
Description=Starts nightly backup

[Timer]
OnCalendar={{ kopia_backup_schedule }}
Unit=kopia.service

[Install]
WantedBy=timers.target
